# etc/supervisor/celeryd_atms.conf
[program:atms_worker]
; Set full path to celery program if using virtualenv
command=/home/raoxin/venv/bin/celery worker -A rtms --loglevel=INFO -n RTMS_worker -Q rtms_task -c 10


; Alternatively,
;command=celery --app=your_app.celery:app worker --loglevel=INFO -n worker.%%h
; Or run a script
;command=celery.sh

directory=/home/raoxin/1111/rtms/rtms
user=root
numprocs=1
stdout_logfile=/var/log/celery/rtms_worker.log
stderr_logfile=/var/log/celery/rtms_worker.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000
